# files/etc/uci-defaults/97_wired_priority
#!/bin/sh

/etc/init.d/wired-priority enable
/etc/init.d/wired-priority restart

exit 0
